# cfg-files/etc/cron.d/vault-cfgbak-cron
SHELL=/bin/bash

#* * * * * root echo "Hello world" >> /var/log/cron.log 2>&1
* * * * * root env - $( cat /root/env.txt ) /usr/local/bin/backup-vault-config.sh >> /var/log/cron.log 2>&1

